# etc/init.d/dnscrypt-proxy
#!/sbin/openrc-run

description="DNSCrypt client proxy daemon"
command="/usr/sbin/dnscrypt-proxy"
command_args="-config /etc/dnscrypt-proxy/dnscrypt-proxy.toml"
command_user="_dnscrypt-proxy"
supervisor=supervise-daemon

depend() {
	need net
	after localmount
	provide dnscrypt-proxy
	use logger
}

start_pre() {
	if ! getent passwd "${command_user}" >/dev/null; then
		eerror "User '${command_user}' does not exist. Please create it first."
		return 1
	fi

	checkpath --directory --owner "${command_user}" --mode 0755 \
		/var/cache/dnscrypt-proxy \
		/var/log/dnscrypt-proxy \
		/run/dnscrypt-proxy
}
