# etc/local.d/dnscrypt_proxy.start
#!/bin/sh

# kill dnsmasq process to free port 53 and retsart dnscrypt-proxy service
killall -9 dnsmasq
rc-service dnscrypt-proxy restart
